# debian/plug-ugmm.cron.daily
#!/bin/sh

php /usr/share/cassa-ugmm/www/expiredmembers.php
